# components/vault-ami/supervisor/run-vaultsetup.conf
[program:vault-setup]
command=/home/ec2-user/vaultsetup.sh
user=ec2-user
stdout_logfile=/home/ec2-user/vault-setup-stdout.log 
stderr_logfile=/home/ec2-user/vault-setup-stderr.log


